# roles/plex/files/plex-update.service
[Unit]
Description=Plex update

[Service]
Type=oneshot
ExecStart=/usr/bin/docker restart plex