# roles/flexget/files/flexget-update.service
[Unit]
Description=FlexGet update

[Service]
Type=oneshot
ExecStart=/usr/bin/docker restart flexget